# docker/supervisord.conf
[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/run/supervisord.pid
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

# ── PostgreSQL (priority 10 — starts first) ───────────────────────────────────
[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/data
user=postgres
priority=10
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/var/log/stockpeak/postgres.log
stderr_logfile=/var/log/stockpeak/postgres.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# ── Next.js (priority 20 — waits 5s for postgres) ─────────────────────────────
[program:nextjs]
command=/bin/bash -c "set -a; source /run/stockpeak.env; set +a && sleep 5 && exec node /app/server.js"
directory=/app
priority=20
autostart=true
autorestart=true
startretries=10
startsecs=3
stdout_logfile=/var/log/stockpeak/nextjs.log
stderr_logfile=/var/log/stockpeak/nextjs.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# ── Python Notifier / APScheduler (priority 30) ────────────────────────────────
[program:notifier]
command=/bin/bash -c "set -a; source /run/stockpeak.env; set +a && sleep 8 && exec python3 /app/scripts/notification_engine.py"
directory=/app
priority=30
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/var/log/stockpeak/notifier.log
stderr_logfile=/var/log/stockpeak/notifier.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
